# services/usb2tcp.service

[Unit]
Description=Broadcast solutions to TCP port

[Service]
ExecStart=/bin/bash /usr/local/bin/usb2tcp.sh

[Install]
WantedBy=multi-user.target

